# market-logger.service.example
[Unit]
Description=Market Logger Scraper
After=network-online.target
Wants=network-online.target
[Service]
Type=simple
User=your-user
WorkingDirectory=/path/to/market-logger
ExecStart=/bin/bash -lc 'source venv/bin/activate && python3 kraken_scraper.py'
Restart=always
RestartSec=5
[Install]
WantedBy=multi-user.target
